# speedydeploy/templates/supervisor/default.conf
{% if supervisor.use_gunicron %}
[program:{{user}}_gunicorn]
directory={{project_path}}/
command={{remote_dir}}/env/bin/gunicorn_django {{project_path}}/settings.py -c {{remote_dir}}/etc/gunicorn/conf.py
user={{user}}
autostart=true
autorestart=true
redirect_stderr=True
{% endif %}

{% if supervisor.use_celery %}
[program:{{user}}_celery]
directory={{project_path}}/
command={{remote_dir}}/env/bin/python {{project_path}}/manage.py celeryd --loglevel=INFO
user={{user}}
numprocs=1
autostart=true
autorestart=true
startsecs=10
redirect_stderr=true
stdout_logfile={{remote_dir}}/log/celeryd.log

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600

; if rabbitmq is supervised, set its priority higher
; so it starts first
priority=998
{% endif %}

{% if supervisor.use_sphinxsearch %}
[program:{{user}}_sphinxsearch]
directory={{project_path}}/
command={{remote_dir}}/searchd
user={{user}}
autostart=true
autorestart=true
redirect_stderr=True
{% endif %}

